# imessage-gateway/com.notif-fwd.imessage-gateway.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.notif-fwd.imessage-gateway</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/path/to/notification-forwarder/imessage-gateway/server.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/notification-forwarder/imessage-gateway</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>IMESSAGE_GATEWAY_PORT</key>
        <string>8095</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/imessage-gateway.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/imessage-gateway.log</string>
</dict>
</plist>
